# recipes-armaz/ertza/ertza/ertza.service
[Unit]
Description=Ertza server
After=armaz-dto.service

[Service]
User=ertza
Group=root
ExecStartPre=-/usr/bin/ertza-ip.sh
ExecStart=/usr/bin/ertza
KillMode=process

[Install]
WantedBy=multi-user.target
