# recipes-armaz/armaz-dto/files/armaz-dto.service
[Unit]
Description=Armaz device tree loader

[Service]
Type=oneshot
ExecStart=/usr/bin/armaz-dto.sh

[Install]
WantedBy=multi-user.target
