# recipes-elias/ertza/ertza/ertza.service
[Unit]
Description=Ertza server

[Service]
User=root
Group=root
ExecStart=/usr/bin/ertza
KillMode=process

[Install]
WantedBy=multi-user.target
